# server/rbtimgdev/systemd/system/pytor.service
[Unit]
Description=Python Tornado Server

[Service]
Type=simple
ExecStart=/usr/bin/pytor

[Install]
WantedBy=multi-user.target
